# meta-pagespeak/recipes-kernel/pagespeak-btn/files/pagespeak-btn.service
[Unit]
Description=Load pagespeak-btn GPIO button driver
# Wait until the RP1 GPIO controller has been probed and /dev/gpiochip0 exists.
# On RPi5, GPIO is on the RP1 chip (PCIe) which probes after modules-load.d runs,
# causing gpio_request_one to return -EPROBE_DEFER if loaded too early.
After=dev-gpiochip0.device
Wants=dev-gpiochip0.device

[Service]
Type=oneshot
ExecStart=/sbin/modprobe pagespeak-btn
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
